# squashfs-root/usr/lib/opkg/info/uhttpd.list
/opt/lib/upgrade/keep.d/uhttpd
/opt/usr/sbin/uhttpd
/opt/etc/init.d/uhttpd
/opt/etc/config/uhttpd
