# squashfs-root/usr/lib/opkg/info/unfs3.list
/etc/exports
/usr/sbin/unfsd
/etc/init.d/unfs3
